# extensions/agent-experience/units/experience-consolidate.service
[Unit]
Description=Agent Experience manual consolidation (package template)
Documentation=https://github.com/misunders2d/pi-experiences#readme
ConditionPathExists=%h/.agents/experience/agent-experience.toml

[Service]
Type=oneshot
# Template only. 0.1.9 does not provide a live consolidation model adapter or package-owned timer.
# Do not enable this unit as-is. Replace ExecStart with an explicit reviewed adapter/fixture command.
ExecStart=/usr/bin/env sh -c 'echo "Agent Experience timer template disabled: replace ExecStart with an explicit reviewed consolidation adapter/fixture command before enabling." >&2; exit 1'
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.agents/experience
