# extensions/agent-experience/units/experience-consolidate.service
[Unit]
Description=Agent Experience manual consolidation (package template)
Documentation=https://github.com/misunders2d/pi-experiences#readme
ConditionPathExists=%h/.agents/experience/agent-experience.toml

[Service]
Type=oneshot
# Template only. Normal manual learning runs through /experience setup analyze-now.
# This package does not provide a package-owned timer. Do not enable this unit as-is.
ExecStart=/usr/bin/env sh -c 'echo "Agent Experience timer template disabled: use /experience setup analyze-now for manual learning; do not enable this timer without an explicit reviewed scheduler design." >&2; exit 1'
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.agents/experience
