# extensions/agent-experience/units/experience-consolidate.timer
[Unit]
Description=Pi Experiences Agent Experience daily scheduled Analyze (setup-rendered template)

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
Unit=pi-experiences-agent-experience-analyze.service

[Install]
WantedBy=timers.target
