# extensions/agent-experience/units/experience-consolidate.timer
[Unit]
Description=Agent Experience daily consolidation timer (package template; manual install only)

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=20m
# Optional laptop safety: uncomment if desired after manual review.
# ConditionACPower=true

[Install]
WantedBy=timers.target
